# auto-pacman-upgrade.service
[Unit]
Description=Auto pacman upgrade service

[Service]
ExecStart=/home/bdsword/Documents/Scripts/AutoPacmanUpgrade/upgrade.sh
Type=forking

